# meta/wip/py-jabber-yahoo-transport/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='pkgsrc/@SMF_NAME@' type='service' version='0'>
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>
    <dependency name='network-service' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/service' />
    </dependency>
    <dependency name='system-log' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/system-log' />
    </dependency>
    <method_context working_directory='@SMF_HOME@'>
      <method_credential user='@SMF_USER@' group='@SMF_GROUP@' />
    </method_context>

    <instance name='@SMF_INSTANCE@' enabled='false'>
      <dependency name='config-files' grouping='require_any' restart_on='refresh' type='path'>
        <service_fmri value='file://localhost@PKG_SYSCONFDIR@/yahoo-transport.xml' />
      </dependency>
      <dependency name='@SMF_INSTANCE@-ejabberd' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/pkgsrc/ejabberd' />
      </dependency>
      <exec_method name='start' type='method' exec='@PREFIX@/libexec/jabber-transports/yahoo/yahoo.py' timeout_seconds='60' />
      <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60' />
      <property_group name='startd' type='framework'>
	<propval name="duration" type="astring" value="child" />
        <propval name='ignore_error' type='astring' value='core,signal' />
      </property_group>
      <template>
        <common_name>
          <loctext xml:lang='C'>PyYahoo-t - the XMPP Yahoo transport</loctext>
        </common_name>
      </template>
    </instance>

    <stability value='Unstable' />
  </service>
</service_bundle>
