# meta/joyent/quickbackup-percona/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='pkgsrc/quickbackup-percona' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <method_context project=':default' resource_pool=':default' working_directory='@VARBASE@/mysql'>
      <method_credential group='mysql' limit_privileges=':default' privileges=':default' supp_groups=':default' user='mysql'/>
      <method_environment>
        <envvar name='PATH' value='@PREFIX@/sbin:@PREFIX@/bin:/usr/sbin:/usr/bin'/>
      </method_environment>
    </method_context>
    <exec_method name='start' type='method' exec='@PREFIX@/bin/quickbackup-percona start' timeout_seconds='18446744073709551615'/>
    <exec_method name='stop' type='method' exec='@PREFIX@/bin/quickbackup-percona stop' timeout_seconds='18446744073709551615'/>
    <exec_method name='restart' type='method' exec='@PREFIX@/bin/quickbackup-percona restart' timeout_seconds='18446744073709551615'/>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='transient'/>
    </property_group>
    <property_group name='quickbackup' type='application'>
      <propval name='backupdir' type='astring' value='@VARBASE@/backups/percona'/>
      <propval name='stream' type='astring' value='xbstream'/>
      <propval name='databases' type='astring' value='all'/>
      <propval name='day' type='astring' value='all'/>
      <propval name='expiredays' type='astring' value='3'/>
      <propval name='hour' type='astring' value='0'/>
      <propval name='minute' type='astring' value='0'/>
      <propval name='username' type='astring' value='username'/>
      <propval name='password' type='astring' value='password'/>
    </property_group>
    <template>
      <common_name>
        <loctext xml:lang='C'>Joyent QuickBackup service for Percona Server (MySQL)</loctext>
      </common_name>
      <documentation>
        <doc_link name='wiki.joyent.com' uri='http://wiki.joyent.com/jpc2/Joyent+Percona+SmartMachine'/>
      </documentation>
    </template>
  </service>
</service_bundle>
